# deploy/spotify-sync.timer
[Unit]
Description=Weekly Spotify liked songs sync
Requires=spotify-sync.service

[Timer]
OnCalendar=Mon 06:00:00
Persistent=true
RandomizedDelaySec=300
Unit=spotify-sync.service

[Install]
WantedBy=timers.target
